# .docker/etc/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770
chown=root:www-data

[supervisord]
nodaemon=true

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:php-fpm]
command=php-fpm -F

[program:nginx]
command=nginx -g 'daemon off;'